# containers/php-cron/crontab
* * * * *  echo "hello" >> /var/log/cron.log 2>&1
# An empty line is required at the end a cron file.